# overlay/lib/systemd/system/sysupgrade_downloader.service
[Unit]
Description=Mininas Sysupgrade Downloader

[Service]
Type=simple
ExecStart=/usr/local/sbin/mininas_sysupgrade_download

