# overlay/lib/systemd/system/sysupgrade_downloader.timer
[Unit]
Description=Mininas Sysupgrade Downloader Timer

[Timer]
OnBootSec=3min
OnUnitActiveSec=4h

[Install]
WantedBy=timers.target
